# templates/docker-system-prune.service.j2
[Unit]
Description=Remove all unused containers, networks, and images.

[Service]
Type=oneshot
ExecStart=docker system prune -af
