# templates/docker-system-prune.timer.j2
[Unit]
Description=Remove all unused containers, networks, and images.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
